# config/services/user/infoplay-webapp.service
[Unit]
Description=InfoPlay webapp server

[Service]
Type=simple
ExecStart=runInfoplayWebapp
Restart=always

[Install]
WantedBy=default.target
