# config/services/user/infoplay-backend.service
[Unit]
Description=InfoPlay backend

[Service]
Type=simple
ExecStart=runInfoplayBackend
Restart=always

[Install]
WantedBy=default.target
